# x1-carbon-drivers/files/trackpoint/fix_trackpoint@.service
[Unit]
Description=Reconnect trackpoint after wakeup
After=sleep.target

[Service]
User=%I
Type=simple
ExecStart=/usr/local/bin/fix_trackpoint
